# freifunk-halle/files/etc/uci-defaults/uci-defaults-freifunk-halle
#!/bin/sh

uci set luci_splash.general.leasetime=5
uci commit luci_splash
